# base/services/syslog-ng.ini
[program:syslog-ng]
command=/usr/sbin/syslog-ng -f /etc/syslog-ng/syslog-ng.conf --foreground
