# haproxy/services/haproxy.ini
[program:haproxy]
command=/usr/sbin/haproxy -V -db -f /etc/haproxy/haproxy.cfg
